# deployment/home-hub-ambient.service
# Home Hub ambient noise monitor — systemd user unit for the Latitude
#
# Waits for the backend to be healthy before starting. Uses the laptop's
# built-in mic (not the Blue Yeti, which is on the Windows desktop).
#
# Install:
#   cp deployment/home-hub-ambient.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable home-hub-ambient.service

[Unit]
Description=Home Hub Ambient Noise Monitor
After=home-hub.service
Wants=home-hub.service

[Service]
Type=simple
WorkingDirectory=%h/home-hub
ExecStartPre=/bin/bash -c 'until curl -sf http://localhost:8000/health; do sleep 2; done'
ExecStart=%h/home-hub/venv/bin/python -m backend.services.pc_agent.ambient_monitor --server http://localhost:8000
Restart=on-failure
RestartSec=10
Environment="PATH=%h/home-hub/venv/bin:/usr/bin"

[Install]
WantedBy=default.target
